# etc/supervisord.d/scylla-server.ini
[program:scylla-server]
directory=/etc/scylla/
command=/opt/qnib/scylla/server/bin/start.sh
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true
