# etc/supervisord.d/scylla-jmx.ini
[program:scylla-jmx]
command=/opt/qnib/scylla/jmx/bin/start.sh
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true
